# crontab
*/15 * * * * imapfilter -v > /var/log/imapfilter.log

